# historian.service
[Unit]
Description=Historian service
After=redis-server.service

[Service]
WorkingDirectory=/opt/termometer
ExecStart=python /opt/termometer/bin/historian.py --config /opt/termometer/etc/config.yaml
Type=simple
User=root
StandardOutput=syslog
StandardError=syslog
Restart=on-failure

[Install]
WantedBy=default.target
